# docker/php/conf/supervisor.d/order.conf
[program:order]
process_name = %(program_name)s_%(process_num)02d
command = php /app/entrypoint/console.php swoole:serve
autostart = false
autorestart = true
user = www-data
numprocs = 1
